# systemd/manga_dl.timer
[Unit]
Description=Execution timer for manga-dl service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target