# systemd/manhua-dl.timer
[Unit]
Description=Execution timer for Manhua-dl service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target